# contrib/systemd/tufmanager.service.in
[Unit]
Description=TUF Manager Server

[Service]
Type=forking
ExecStart=@LIBEXECDIR@/tuf-server

[Install]
WantedBy=multi-user.target
