# data/systemd/tufmanager.service.in
[Unit]
Description=TUF Manager Server

[Service]
Type=forking
PIDFile=/run/tufmanager/pid
ExecStart=@LIBEXECDIR@/tuf-server

[Install]
WantedBy=multi-user.target
